# crontab
PATH=/usr/local/bin:/usr/local/sbin:/usr/bin:/usr/sbin:/bin:/sbin

# 1. Sync Omramin (Omron data to Garmin) at 08:10 AM
10 8 * * * omramin sync > /proc/1/fd/1 2>&1

# 2. Daily Sync GarminGo: Fetch Yesterday (final) + Today (sleep) at 08:15 AM
15 8 * * * cd /app && python3 src/main.py cli-sync --profile USER1 --start-date $(date -d "yesterday" +\%Y-\%m-\%d) --end-date $(date +\%Y-\%m-\%d) > /proc/1/fd/1 2>&1
