# crontab
PATH=/usr/local/bin:/usr/local/sbin:/usr/bin:/usr/sbin:/bin:/sbin

# Daily Sync at 08:10 AM: Omron BP first, then Garmin to Google Sheets
10 8 * * * cd /app && omramin --keyring-backend file --keyring-file /root/.config/omramin/config.tokens.json sync --days 1 && python3 -m src.main cli-sync --start-date $(date +\%Y-\%m-\%d) --end-date $(date +\%Y-\%m-\%d) --profile USER1 --output-type sheets > /proc/1/fd/1 2>&1
